# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx 

[program:flask]
command=/usr/bin/python orbital.py 2&1> /orbital/webapp/output.log
directory=/orbital/webapp
